# source/_posts/Prometheus/alertmanager/alertmanager.ini
[program:alertmanager]
command=/opt/fonsview/3RD/alertmanager/alertmanager 
directory=/opt/fonsview/3RD/alertmanager
autostart=true
startsecs=3
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups = 10
stdout_logfile=/opt/fonsview/3RD/alertmanager/alertmanager.log
stopasgroup=false
killasgroup=false